# data/indicator-application.service.in
[Unit]
Description=Indicator Application Service
PartOf=graphical-session.target

[Service]
ExecStart=@pkglibexecdir@/indicator-application-service
