# data/indicator-application.service.in
[Unit]
Description=Indicator Application Service
PartOf=graphical-session.target
After=unity-panel-service unity8

[Service]
ExecStart=@pkglibexecdir@/indicator-application-service
